# scripts/com.pyoung.actress-profiles.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pyoung.actress-profiles</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/pyoung/workspace/organizer3/scripts/run-actress-profiles.sh</string>
  </array>

  <!-- Fire at 2:00 AM daily -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/pyoung/workspace/organizer3/logs/actress-profile-launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/pyoung/workspace/organizer3/logs/actress-profile-launchd.log</string>

  <!-- Run immediately on load if the last scheduled fire was missed -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
